# service/openrc/mcp
#!/sbin/runscript

start() {
	ebegin "Starting MCP"
	start-stop-daemon --start --quiet --exec /usr/bin/mcp --background --pidfile /var/run/mcp.pid --make-pidfile
	eend $?
}

stop() {
	ebegin "Stopping MCP"
	start-stop-daemon --stop --quiet --pidfile /var/run/mcp.pid
	eend $?
}
