# dist/init/openrc/armaadmin
#!/sbin/runscript

start() {
	ebegin "Starting ArmaAdmin"
	start-stop-daemon --start --quiet --exec "`which armaadmin`" --background --pidfile /var/run/armaadmin.pid --make-pidfile
	eend $?
}

stop() {
	ebegin "Stopping ArmaAdmin"
	start-stop-daemon --stop --quiet --pidfile /var/run/armaadmin.pid
	eend $?
}
