# dist/init/openrc/armaadmin
#!/sbin/runscript

start() {
	ebegin "Starting ArmaAdmin"
	start-stop-daemon --start --exec /usr/bin/armaadmin --pidfile /var/run/armaadmin.pid --make-pidfile
	eend $?
}
